# airootfs/etc/skel/.config/autostart/nm-applet.desktop
[Desktop Entry]
Type=Application
Name=Network Manager Applet
Exec=nm-applet
X-GNOME-Autostart-enabled=false
NoDisplay=true
