# airootfs/etc/systemd/system/greetd.service.d/10-plymouth.conf
[Unit]
After=plymouth-start.service
Wants=plymouth-quit.service

[Service]
# Tell Plymouth to hand off as soon as greetd is up, avoiding the 60s timeout
ExecStartPost=-/usr/bin/plymouth deactivate
ExecStopPost=-/usr/bin/plymouth quit --wait
